# rpm/files/systemd/micro-platform-user-api.service
[Unit]
Description=Micro Platform - User API
After=micro-platform-infra.target
Requires=micro-platform-infra.target
PartOf=micro-platform-all.target

[Service]
Type=exec
User=nobody
Group=nobody
WorkingDirectory=/opt/micro-platform/apps/user-api

# Load shared configurations first
EnvironmentFile=/opt/micro-platform/apps/conf-shared/env.properties
EnvironmentFile=/opt/micro-platform/apps/conf-shared/redis.properties

# Load app-specific configuration (can override shared)
EnvironmentFile=/opt/micro-platform/apps/user-api/user-api.properties

# Optional override file (won't fail if missing)
EnvironmentFile=-/opt/micro-platform/apps/user-api/user-api.override

ExecStart=/opt/micro-platform/apps/user-api/user-api
Restart=always
RestartSec=10
StandardOutput=append:/var/log/micro-platform/user-api/stdout.log
StandardError=append:/var/log/micro-platform/user-api/stderr.log
# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/log/micro-platform/user-api

[Install]
WantedBy=multi-user.target

